# root/etc/services.d/openvpn/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc openvpn \
	--config /config/openvpn.ovpn \
	--log /config/log/openvpn.log
